# crontab
* * * * * root python /opt/geotechnical-data-platform-backup-utility/backup.py
